# php-supervisor/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

loglevel=error

[program:fpm]
command=/usr/sbin/php-fpm7 -R --nodaemonize

[program:ssh]
command=/usr/sbin/sshd -D

[include]
files = /etc/supervisor/conf.d/*.conf